# clipboard-daemon/cpd.service
[Unit]
Description=Daemon for project clipboard-buffer
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/bin/touch %t/clipboard-buffer/cpdA.txt %t/clipboard-buffer/cpdB.txt
ExecStart=/usr/local/bin/cpd
RuntimeDirectory=clipboard-buffer
TimeoutStopSec=5

[Install]
WantedBy=graphical-session.target